# ops/hermes-gateway.service.d/hermes-env.conf
# Drop-in for gateway-adapter env (sibling of slack-env.conf).
# Leading `-` tolerates an absent file; __TARGET_DIR__ is templated by
# setup-hermes.sh. Operator stages contents via stage-hermes-env.sh.

[Service]
EnvironmentFile=-__TARGET_DIR__/auth/hermes.env
