# ops/hermes-gateway.service.d/hermes-env.conf
# Drop-in layered onto the canonical hermes-gateway.service. Sibling of
# slack-env.conf. The leading `-` on EnvironmentFile= keeps the unit
# loading even when the file is briefly absent (e.g. during ops work).
# __TARGET_DIR__ is substituted by setup-hermes.sh at install time.
#
# auth/hermes.env carries gateway-adapter tokens that aren't Slack-bound.
# LINEAR_API_KEY is the v0 entry — staged via stage-hermes-env.sh and
# consumed by the linear skill when the gateway agent is asked to file
# an issue.

[Service]
EnvironmentFile=-__TARGET_DIR__/auth/hermes.env
