# ops/hermes-gateway.service.d/ops-prod-env.conf
# Drop-in for PROD ops env (sibling of ops-env.conf). Carries
# AWS_PROD_ACCESS_KEY_ID / AWS_PROD_SECRET_ACCESS_KEY used by the
# aws-lambda-debug-prod and dynamodb-query-prod skills.
#
# Kept in a separate file from ops.env so prod credentials are visible at
# the filesystem level. Skills remap AWS_PROD_* onto AWS_ACCESS_KEY_ID /
# AWS_SECRET_ACCESS_KEY per-command, so forgetting the remap falls back
# to the non-prod creds — which is wrong but doesn't leak prod access.
#
# Leading `-` tolerates an absent file; __TARGET_DIR__ is templated by
# setup-hermes.sh. Operator stages contents via stage-secrets.sh.

[Service]
EnvironmentFile=-__TARGET_DIR__/auth/ops-prod.env
