# ops/hermes-gateway.service.d/ops-env.conf
# Drop-in for non-prod ops env (sibling of hermes-env.conf).
# Carries AWS dev/staging/test credentials and NEW_RELIC_API_KEY used by
# the aws-lambda-debug, dynamodb-query, and new-relic-lambda skills.
# Leading `-` tolerates an absent file; __TARGET_DIR__ is templated by
# setup-hermes.sh. Operator stages contents via stage-secrets.sh.

[Service]
EnvironmentFile=-__TARGET_DIR__/auth/ops.env
