# ops/hermes-gateway.service.d/runtime-env.conf
# Drop-in for values-derived env vars (sibling of slack-env.conf and
# hermes-env.conf). The file at __TARGET_DIR__/auth/gateway-runtime.env is
# rewritten on every setup-hermes.sh run from deploy.values.yaml — it
# carries non-secret config that the gateway reads from env (e.g.
# SLACK_ALLOWED_USERS). Do not hand-edit the file; edit deploy.values.yaml
# and re-run the installer.
#
# Leading `-` tolerates an absent file (early-stage installs);
# __TARGET_DIR__ is templated by setup-hermes.sh.

[Service]
EnvironmentFile=-__TARGET_DIR__/auth/gateway-runtime.env
