# ops/hermes-gateway.service.d/slack-env.conf
# Drop-in layered onto the canonical hermes-gateway.service. The leading
# `-` on EnvironmentFile= keeps the unit loading even when a file is
# briefly absent (e.g. during ops work). __TARGET_DIR__ is substituted
# by setup-hermes.sh at install time.

[Service]
EnvironmentFile=-/etc/default/hermes-gateway
EnvironmentFile=-__TARGET_DIR__/auth/slack.env
