# ops/hermes-gateway.service.d/slack-env.conf
# Drop-in for hermes-gateway.service. The canonical unit is generated by
# `hermes gateway install --system` (see hermes_cli/gateway.py) and gets
# refreshed automatically on `hermes update` to track upstream changes.
#
# This drop-in carries the bits that are local to this fork's deploy:
#   - /etc/default/hermes-gateway   non-secret operator config (channel
#                                   allowlist, optional overrides). 0644.
#   - <TARGET>/auth/slack.env       Slack bot/app tokens, hand-staged at
#                                   provisioning time. 0640 root:agent —
#                                   matches the github-app.env pattern.
#
# Both files are referenced with a leading `-` so a missing file degrades
# to "unit loads, gateway runs without those env values" rather than a
# load-time failure. The Slack adapter's auto-enable hinges on
# SLACK_BOT_TOKEN being present in the env at startup.
#
# __TARGET_DIR__ is substituted by setup-hermes.sh at install time so
# non-default --target installs (e.g. /opt/hermes) point at the right path.

[Service]
EnvironmentFile=-/etc/default/hermes-gateway
EnvironmentFile=-__TARGET_DIR__/auth/slack.env
