# ops/atlas-source-sync.service
[Unit]
Description=Atlas source sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__AGENT_USER__
Group=__AGENT_USER__
Environment=HERMES_HOME=__TARGET_DIR__
Environment=ATLAS_CONFIG=__ATLAS_CONFIG__
Environment=ATLAS_SYNC_SOURCE_NAME=__ATLAS_SYNC_SOURCE_NAME__
Environment=HERMES_ATLAS_SOURCE_READER_GH_CONFIG=__TARGET_DIR__/auth/atlas-source-reader.env
EnvironmentFile=-/etc/default/atlas-source-sync
ExecStart=/usr/local/sbin/atlas-source-sync-runner
Restart=no
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
