# ops/com.hermes.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/hermes-sync</string>
    </array>

    <!-- Drop privilege from root (LaunchDaemon) to the agent user before exec. -->
    <key>UserName</key>
    <string>hermes</string>
    <key>GroupName</key>
    <string>hermes</string>

    <!-- 2-minute cadence; override by editing this plist or `launchctl unload`
         and replacing it. AccuracySec / coalescing not supported on launchd. -->
    <key>StartInterval</key>
    <integer>120</integer>

    <!-- Don't keep alive on its own — the timer itself drives invocation. -->
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <true/>

    <!-- Logs land in /var/log so an operator can `tail -f` without chasing
         per-user paths. The plist itself is root-owned so the agent can't
         redirect or silence them. -->
    <key>StandardOutPath</key>
    <string>/var/log/hermes-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/hermes-sync.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
