# ops/hermes-reviewer-token.timer
[Unit]
Description=Refresh /run/hermes/reviewer-gh-token from the reviewer GitHub App

[Timer]
# Mint shortly after boot so the reviewer worker has a valid token by the
# time the first quay-tick fires — the file lives in /run/, which is
# wiped on reboot.
OnBootSec=30s
# 30-min cadence. The helper's internal 5-min refresh margin means the
# on-disk token only actually rotates once per ~55 min of token TTL;
# this cadence keeps both layers comfortably under GitHub's 1-hour cap
# without thrashing the cache.
OnUnitActiveSec=30min
AccuracySec=30s
Persistent=true
Unit=hermes-reviewer-token.service

[Install]
WantedBy=timers.target
