# ops/hermes-gateway.service
[Unit]
Description=Hermes messaging gateway (Slack, etc.)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# __AGENT_USER__ and __TARGET_DIR__ are substituted by setup-hermes.sh at
# install time. The unit file in-repo is a template, not a runnable unit.
User=__AGENT_USER__
Group=__AGENT_USER__
# /etc/default/hermes-gateway carries non-secret config (HERMES_HOME, the
# SLACK_ALLOWED_CHANNELS allowlist, etc). slack.env carries the bot/app
# tokens, mode 0640 root:agent. Both are optional (`-` prefix) so the unit
# can load even when one file is briefly missing during ops work.
EnvironmentFile=-/etc/default/hermes-gateway
EnvironmentFile=-__TARGET_DIR__/auth/slack.env
ExecStart=__TARGET_DIR__/hermes-agent/venv/bin/hermes gateway run --replace
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
