# ops/quay-orchestrator.service
[Unit]
Description=Quay orchestrator delivery outbox drain
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__AGENT_USER__
Group=__AGENT_USER__
Environment=HERMES_HOME=__TARGET_DIR__
Environment=QUAY_DATA_DIR=__TARGET_DIR__/quay
Environment=QUAY_ORCHESTRATOR_CONFIG=__TARGET_DIR__/quay/orchestrator.json
EnvironmentFile=-__TARGET_DIR__/auth/hermes.env
EnvironmentFile=-__TARGET_DIR__/auth/quay.env
EnvironmentFile=-__TARGET_DIR__/auth/slack.env
EnvironmentFile=-/etc/default/quay-orchestrator
ExecStart=/usr/local/sbin/quay-orchestrator-runner
Restart=no
# The deployed runner parks human waits and exits; later timer ticks poll
# waiting_human handoffs without keeping this oneshot active for reply_timeout.
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
