# ops/quay-tick.service
[Unit]
Description=Quay tick — drive queued tasks forward
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# A tick spawns a detached tmux server (`tmux new-session -d`) for each
# worker. systemd's default KillMode=control-group SIGKILLs every
# process in the unit's cgroup when ExecStart returns — including the
# tmux server we just forked off, before pipe-pane has flushed the first
# byte and before the wrapper writes .quay-exit-code. Limiting kill
# scope to the ExecStart process leaves detached children (tmux server
# + pane shell + claude) alone so workers actually run to completion.
KillMode=process
# User/Group and the target-dir-derived paths below are substituted by
# setup-hermes.sh at install time, so non-default --user / --target
# values get a correctly-targeted unit (matches hermes-sync).
User=__AGENT_USER__
Group=__AGENT_USER__
# QUAY_DATA_DIR pins the data dir (sqlite, worktrees, repo bare clones,
# logs) — set here rather than in config.toml so the runtime path lives
# in exactly one place that the installer can reason about. HERMES_HOME
# tells the GH-token helper (invoked by quay-tick-runner before each
# tick) which install root to read auth/github-app.{pem,env} from.
Environment=QUAY_DATA_DIR=__TARGET_DIR__/quay
Environment=HERMES_HOME=__TARGET_DIR__
# Reviewer App credentials live outside HERMES_HOME because they are a
# separate deploy-rendered identity. quay-tick-runner reads this file and
# exports QUAY_REVIEWER_GH_TOKEN for reviewer attempts.
Environment=HERMES_REVIEWER_GH_CONFIG=/etc/hermes/reviewer.env
# The reviewer helper cache points at /run/hermes/reviewer-token-cache.json.
# Let systemd create that directory as the agent user before each tick.
RuntimeDirectory=hermes
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
# Adapter tokens (LINEAR_API_KEY etc.) live in auth/quay.env, written by
# stage-secrets.sh out-of-band. Leading `-` makes the file optional so
# the unit starts even before tokens are staged: a tick without a Linear
# adapter just has nothing to enqueue.
EnvironmentFile=-__TARGET_DIR__/auth/quay.env
# Operator overrides (cadence retunes, ad-hoc env vars). Preserved
# across installer re-runs — same /etc/default/ pattern as hermes-sync.
EnvironmentFile=-/etc/default/quay-tick
# quay-tick-runner mints $GH_TOKEN and $QUAY_REVIEWER_GH_TOKEN from the
# installer's GitHub App helper, then execs `quay tick`. Without this,
# the worker pipeline's `gh pr list/create` calls fail for lack of HTTPS
# auth — the deploy keys cover git push, not the API — and reviewer
# attempts cannot post as the reviewer App.
ExecStart=/usr/local/sbin/quay-tick-runner
# A single bad tick should not lock the timer out forever. Restart=no
# keeps the systemd state honest — failures show up in `systemctl
# list-units`.
Restart=no
# Generous: a tick may spawn a worker tmux session; the supervisor lock
# in quay prevents overlapping ticks regardless.
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
